# pigall/02-setup-pigall/files/flask.service
[Unit]
Description = Start pigall web server.

[Service]
ExecStart = /home/pi/webapp/start-webapp.sh

[Install]
WantedBy = multi-user.target
